# agent/crontab
* * * * * test -f /usr/local/src/halt.txt || go run /usr/local/src